# deploy/supervisor/flaskapp.conf
[program:flaskapp]
directory=/code
command=/usr/bin/python ./app.py
startsecs=10
stopwaitsecs=0
autostart=true
autorestart=true
stopsignal=INT
stopasgroup=true
killasgroup=true
